# libdecor_lsm/libdecor_lsm.service
[Unit]
Description=Block libdecor GTK plugin

[Service]
Type=notify
ExecStart=%Y/../../../libexec/libdecor_lsm

[Install]
WantedBy=graphical.target
